# .provision/roles/autohome/files/supervisor/cmdb.conf
[program:cmdb]
command=/bin/sh /home/www/cmdb_shells/cmdb.sh

stdout_logfile = /var/log/uwsgi_access.log